# LaunchDaemon/com.youtubeminus.helper.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<!-- Launchd label — must match com.youtubeminus.helper exactly -->
	<key>Label</key>
	<string>com.youtubeminus.helper</string>

	<!-- Path where the install script copies the binary -->
	<key>Program</key>
	<string>/Library/PrivilegedHelperTools/YouTubeMinusHelper</string>

	<!-- Start at boot, before login session -->
	<key>RunAtLoad</key>
	<true/>

	<!-- Restart automatically if it crashes -->
	<key>KeepAlive</key>
	<true/>

	<!-- Register a Mach service so the main app can connect via XPC -->
	<key>MachServices</key>
	<dict>
		<key>com.youtubeminus.helper</key>
		<true/>
	</dict>

	<!-- Capture logs to a persistent file -->
	<key>StandardOutPath</key>
	<string>/Library/Logs/YouTubeMinusHelper.log</string>
	<key>StandardErrorPath</key>
	<string>/Library/Logs/YouTubeMinusHelper.log</string>
</dict>
</plist>
